# udev.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/home/pi/scripts/retropie-standby/wakeup.sh"
